# files/etc/hotplug.d/iface/01-restart-socket-xp
logger -t restart_socketxp "Network interface changed. Restarting SocketXP"

service socket_xp restart

logger -t restart_socketxp "SocketXP restarted"
